# x11/etc/supervisord.d/websockify.conf
[program:websockify]
command=websockify 127.0.0.1:9900 127.0.0.1:5900
autorestart=true
stdout_logfile=/var/log/supervisor/websockify.log
redirect_stderr=true
user=nobody
