# base/etc/profile.d/locale.sh
# Use built-in C locale rather than an English variant, but with UTF-8
: ${LANG:=C.UTF-8}; export LANG
