# base/etc/profile.d/prompt.sh
export PS1='\[\033[01;32m\]\W\[\033[0m\] \$ '
eval `dircolors -b`
alias ls='ls --color=auto' 2>/dev/null
